# macos/LaunchAgents/com.haacked.sprint-status.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.haacked.sprint-status</string>

    <!--
        Use zsh login shell so .zshenv provides PATH (Homebrew, ~/.local/bin,
        Cargo) and path_helper picks up entries from /etc/paths.d.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string><![CDATA[
log_dir="$HOME/.local/state/sprint-status"
mkdir -p "$log_dir"
exec "$HOME/.dotfiles/bin/sprint-status-run" >> "$log_dir/launchd.log" 2>&1
]]></string>
    </array>

    <!--
        Fire Friday at 08:00 local time, posted to #team-feature-flags as the
        end-of-sprint-week status. StartCalendarInterval uses the system's
        local time zone (America/Los_Angeles on this machine, i.e. PT);
        Weekday 5=Friday. 15 min before standup's Friday 08:15 run so the two
        headless Claude runs don't overlap.
    -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
</dict>
</plist>
