# macos/LaunchAgents/com.haacked.secretive-ssh-auth.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.haacked.secretive-ssh-auth</string>
    <!--
    Points the session-wide SSH_AUTH_SOCK default at the stable agent
    symlink (see zsh/zshrc.symlink) rather than Secretive's socket
    directly, so processes that never source .zshrc (GUI apps, tools
    that spawn non-interactive shells) still pick up a forwarded agent
    when one is live. Seeds the symlink to the local Secretive socket
    first in case this runs before any shell has initialized it. Also
    seeds agent.pub from user.localSigningKey so GUI apps that sign a
    commit before any interactive shell has run don't find it missing.
    Bare `git` resolves fine here: launchd runs jobs with PATH=/usr/bin:
    /bin:/usr/sbin:/sbin, and /usr/bin/git (the Xcode CLT shim, a
    Homebrew prerequisite already assumed elsewhere in this repo) is on
    that PATH.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>mkdir -p "$HOME/.ssh" &amp;&amp; ln -sf "$HOME/Library/Containers/com.maxgoedjen.Secretive.SecretAgent/Data/socket.ssh" "$HOME/.ssh/agent.sock" &amp;&amp; /bin/launchctl setenv SSH_AUTH_SOCK "$HOME/.ssh/agent.sock"; LOCAL_SIGNING_KEY=$(git config --global --includes --get user.localSigningKey 2>/dev/null); [ -n "$LOCAL_SIGNING_KEY" ] &amp;&amp; [ -f "$LOCAL_SIGNING_KEY" ] &amp;&amp; cat "$LOCAL_SIGNING_KEY" &gt; "$HOME/.ssh/agent.pub"</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
